# qr3_79.conf
[program:qr3_79]
command=python3 /home/curator3/qr3_79/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_79/test.err.log
stdout_logfile=/home/curator3/qr3_79/test.err.log
